# house-keeping/snapshot-cleaner.service
[Unit]
Description=Etcd snapshot cleaner
[Service]
ExecStart=/usr/bin/find /etc/etcd/backup/ -name "*.*" | xargs rm
